# launchd/com.paz.moki.telegram-poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
                      "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.paz.moki.telegram-poll</string>

    <!-- Poll Telegram every 5 minutes (300 seconds) for button clicks.
         Independent of pipeline / autonomy so approvals feel instant. -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
        <string>/Users/ASUS/education-agents/telegram_approval.py</string>
        <string>--poll</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ASUS/education-agents</string>

    <key>StandardOutPath</key>
    <string>/Users/ASUS/education-agents/output/telegram_poll.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ASUS/education-agents/output/telegram_poll.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
